# pacman-optimize.service
[Unit]
Description=pacman optimize database

[Service]
ExecStart=/usr/bin/pacman-optimize

